# src/roles/iop_vex_downloader/files/iop-vex-download.path
[Unit]
Description=Path unit for vex-latest.tar.zst manual file watcher

[Path]
PathChanged=/var/lib/foreman/vex-latest.tar.zst
PathModified=/var/lib/foreman/vex-latest.tar.zst

[Install]
WantedBy=multi-user.target
