# src/roles/iop_vex_downloader/templates/iop-vex-download.timer.j2
[Unit]
Description=Timer for vex-latest.tar.zst download

[Timer]
OnActiveSec=0
OnUnitActiveSec={{ iop_vex_downloader_timer_interval }}

[Install]
WantedBy=timers.target
