# src/roles/foreman/templates/foreman.socket.j2
[Unit]
Description=Foreman socket

[Socket]
ListenStream={{ foreman_listen_stream }}
SocketUser=apache
SocketMode=0600

NoDelay=false
ReusePort=true
Backlog=1024

[Install]
WantedBy=sockets.target
