# src/roles/foreman/templates/foreman.socket.j2
[Unit]
Description=Foreman socket

[Socket]
ListenStream={{ foreman_listen_stream }}

[Install]
WantedBy=sockets.target
