# src/roles/foreman/templates/foreman-db-sessions-clear.timer.j2
[Unit]
Description=Timer for Foreman db:sessions:clear task

[Timer]
OnCalendar=*-*-* 23:15:00
Persistent=true

[Install]
WantedBy=timers.target
