# src/roles/foreman/templates/foreman-reports-daily.timer.j2
[Unit]
Description=Timer for Foreman reports:daily task

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
